# configs/supervisord_conf.d/supervisord.conf
[program:funs] 
command=uwsgi /etc/uwsgi/conf.d/config.ini
directory=/opt/project
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/project/uwsgi_supervisor.log
